# images/aredn-base/rootfs/etc/s6/babeld/run
#!/bin/sh

exec babeld -c /etc/babel.conf -c /tmp/babel-generated.conf