# images/mesh-base/rootfs/etc/s6/olsrd/run
#!/bin/sh

exec olsrd -f /etc/olsrd/olsrd.conf -pidfile /tmp/olsrd.pid -nofork