# images/weewx/rootfs/etc/s6/weewx/run
#!/bin/sh

exec su weewx bash -c '/usr/local/bin/weewxd --config ${CONFIG_OVERRIDE:-/etc/weewx/weewx.conf}'
