# images/weewx/rootfs/etc/s6/weewx/run

#!/bin/sh

exec doas -u weewx -- \
    /usr/bin/weewxd --config ${CONFIG_OVERRIDE:-/etc/weewx/weewx.conf}
