# images/aredn-base/rootfs/etc/s6/vtund/run
#!/bin/sh

exec vtund -s -f /etc/vtundsrv.conf -n