# images/mesh-base/rootfs/etc/s6/messages/run
#!/bin/sh

exec tail -F /var/log/messages /var/log/cron.log 2>/dev/null